# roles/cronn/files/crontab
@hourly /scripts/backup-data.sh
@hourly /scripts/duck.sh
